# extras/fools2019.conf
[program:fools2019]
directory=/srv/fools2019/server
command=python3 /srv/fools2019/server/main.py
user=fools2019
environment=HOME="/srv/fools2019",USER="fools2019"
autostart=true
autorestart=true
stdout_logfile=/var/log/fools2019.stderr.log